# pi-unattended-upgrades/reboot-check.service
[Unit]
Description=Check for a pending reboot (unattended-upgrades) and notify via Apprise

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/notify-reboot-required.sh
